# systemd/max-identity-bot.service
[Unit]
Description=MAX Identity Bot safe test service
Wants=network-online.target
After=network-online.target

[Service]
Type=simple
WorkingDirectory=/root/zabbix-max-alert-bot
EnvironmentFile=-/root/zabbix-max-alert-bot/.env
ExecStart=/root/.nvm/versions/node/v24.16.0/bin/node /root/zabbix-max-alert-bot/src/bot-platform/app.js
Restart=always
RestartSec=5
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
